# oadcgs-es-elastic/scripts/crons/serena_cleanup.cron
#!/bin/bash

SCRIPT_DIR="/etc/logstash/scripts/"
/etc/logstash/scripts/venv/bin/python ${SCRIPT_DIR}SerenaCleanup.py
